# webot/blog/media/configs/21sites/bbs.tianya.cn/tianya_bbs.cron
0 */1 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=3600 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-00.conf
1 */1 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=3600 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-01.conf
2 */1 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=3600 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-02.conf
3 */1 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=3600 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-03.conf
4 */1 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=3600 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-04.conf
5 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-05.conf
6 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-06.conf
7 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-07.conf
8 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-08.conf
9 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-09.conf
10 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-10.conf
11 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-11.conf
12 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-12.conf
13 */2 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=7200 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-13.conf
14 */3 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=10800 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-14.conf
15 */3 * * * curl http://192.168.3.173:6800/schedule.json -d project=example -d spider=example -d setting=CLOSESPIDER_TIMEOUT=10800 -d setting=DEPTH_LIMIT=5 -d config=http://192.168.3.175/bbs/tianya/tianya_bbs-15.conf
